# bot.conf
[program:oblacko]
directory=/home/ubuntu/tokiobae_bot
command=python3 app.py

autostart=true
autorestart=true
environment=HOME="/home/ubuntu",USER="ubuntu"

stderr_logfile=/home/ubuntu/tokiobae_bot/logfile_err.log
stdout_logfile=/home/ubuntu/tokiobae_bot/logfile.log
