# testing_bot.conf
[program:testing_bot]
directory=/home/ubuntu/pharmline_testing_bot
command=python3 app.py

autostart=true
autorestart=true
environment=HOME="/home/ubuntu",USER="ubuntu"

stderr_logfile=/home/ubuntu/pharmline_testing_bot/logfile_err.log
stdout_logfile=/home/ubuntu/pharmline_testing_bot/logfile.log
